# system_files/homeserver/etc/systemd/system/decrypt-sshkey.service
[Unit]
Description=Decrypt SSH private key from SOPS
After=local-fs.target
Wants=network-online.target

[Service]
Type=oneshot
# Decrypt into /var/lib/secrets/ssh, which lives on /var (writable on bootc)
ExecStart=/usr/bin/mkdir -p /var/lib/secrets/ssh
# /usr/share is read-only, world-readable image content; /var/lib is writable, host state
ExecStart=/usr/bin/sops --decrypt /usr/share/secrets/ssh/id_ed25519.sops > /var/lib/secrets/ssh/id_ed25519
Environment=SOPS_AGE_KEY_FILE=/etc/sops/age/keys.txt

[Install]
WantedBy=multi-user.target
