# system_files/homeserver/etc/systemd/system/decrypt-sshkey.service
[Unit]
Description=Decrypt SSH private key from SOPS
After=local-fs.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=SOPS_AGE_KEY_FILE=/etc/sops/age/keys.txt
ExecStart=/usr/bin/install -d -m 0700 /var/lib/secrets/ssh
# /usr/share is read-only image content; /var/lib is writable, host state
ExecStart=/usr/bin/sops --decrypt --output /var/lib/secrets/ssh/id_ed25519 /usr/share/secrets/ssh/id_ed25519.sops
ExecStart=/usr/bin/chmod 600 /var/lib/secrets/ssh/id_ed25519
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target