# system_files/homeserver/etc/profile.d/90-registry-auth.sh
export ENV_LOAD=profiled
export GITHUB_ACTOR=cthtrifork