# system_files/homeserver/etc/profile.d/90-registry-auth.sh
export ENV_LOAD=profiled
export GITHUB_ACTOR=cthtrifork
export AGE_PUB=age1qth3nzg27zdmlzpjl2gedzmt0e5nly84q82knja3fcq365kw6u9qeecy6j
export SOPS_AGE_KEY_FILE=/etc/sops/age/keys.txt