# system_files/homeserver/etc/profile.d/90-registry-auth.sh
ENV_LOAD=profiled
GITHUB_ACTOR=cthtrifork