# system_files/jumphost/usr/share/containers/systemd/demo-backend.container
[Unit]
Description=Demo HTTPS backend (internal.service.local) via Caddy
After=network-online.target
Wants=network-online.target

[Container]
ExecStartPre=/usr/bin/install -d -m 0755 /var/lib/caddy-demo
Image=docker.io/library/caddy:2
PublishPort=8443:8443
Volume=/etc/caddy-demo:/etc/caddy
Volume=/var/lib/caddy-demo:/data


[Service]
Restart=always

[Install]
WantedBy=multi-user.target